# config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper periodically

[Timer]
OnStartupSec=1
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
